# templates/thanos-store.service.erb
# This file managed by Puppet
# Template Path: <%= @module_name %>/templates/thanos-store.service
### Managed by Puppet ###
[Unit]
Description=Prometheus Thanos Subsystem
Wants=basic.target
After=basic.target network.target

[Service]
ExecStart=/usr/bin/thanos store \
<%- if @gcloudtrace_project -%>
  --gcloudtrace.project=<%=@gcloudtrace_project%> \
<%- end -%>
<%- if @gcloudtrace_sample_factor -%>
  --gcloudtrace.sample-factor=<%=@gcloudtrace_sample_factor%> \
<%- end -%>
<%- if @grpc_address -%>
  --grpc-address=<%=@grpc_address%> \
<%- end -%>
<%- if @grpc_advertise_address -%>
  --grpc-advertise-address=<%=@grpc_advertise_address%> \
<%- end -%>
<%- if @http_address -%>
  --http-address=<%=@http_address%> \
<%- end -%>
<%- if @http_advertise_address -%>
  --http-advertise-address=<%=@http_advertise_address%> \
<%- end -%>
<%- if @cluster_address -%>
  --cluster.address=<%=@cluster_address%> \
<%- end -%>
<%- if @cluster_advertise_address -%>
  --cluster.advertise-address=<%=@cluster_advertise_address%> \
<%- end -%>
<%- if @cluster_peers -%>
  <% @cluster_peers.each do |val| -%>
--cluster.peers=<%= val %> \
  <% end -%>
<%- end -%>
<%- if @cluster_gossip_interval -%>
  --cluster.gossip-interval=<%=@cluster_gossip_interval%> \
<% end -%>
<%- if @cluster_pushpull_interval -%>
  --cluster.pushpull-interval=<%=@cluster_pushpull_interval%> \
<%- end -%>
<%- if @cluster_refresh_interval -%>
  --cluster.refresh-interval=<%=@cluster_refresh_interval%> \
<%- end -%>
<%- if @cluster_secret_key -%>
  --cluster.secret-key=<%=@cluster_secret_key%> \
<%- end -%>
<%- if @cluster_network_type -%>
--cluster.network-type=<%=@cluster_network_type%> \
<%- end -%>
<%- if @data_dir -%>
  --data-dir=<%=@data_dir%> \
<%- end -%>
<%- if @store_objstore_config_file -%>
  --objstore.config-file=<%=@store_objstore_config_file%> \
<%- end -%>
<%- if @index_cache_size -%>
  --index-cache-size=<%=@index_cache_size%> \
<%- end -%>
<%- if @chunk_pool_size -%>
  --chunk-pool-size=<%=@chunk_pool_size%> \
<%- end -%>
<%- if @log_level -%>
  --log.level=<%=@log_level%>
<%- else -%>
  --log.level=info
<% end -%>
Restart=always

[Install]
WantedBy=multi-user.target

