# templates/thanos.service.erb
# This file managed by Puppet
# Template Path: <%= @module_name %>/templates/thanos.service
### Managed by Puppet ###
[Unit]
Description=Thanos Sidecar
Wants=basic.target
After=basic.target network.target

[Service]
Environment='S3_SECRET_KEY=<%=@s3_secret_key%>'
Environment="HOSTNAME=%H"
LimitNOFILE=infinity
LimitNPROC=infinity
ExecStart=/usr/bin/thanos sidecar \
<%- if @gcloudtrace_project -%>
  --gcloudtrace.project=<%=@gcloudtrace_project%> \
<%- end -%>
<%- if @gcloudtrace_sample_factor -%>
  --gcloudtrace.sample-factor=<%=@gcloudtrace_sample_factor%> \
<%- end -%>
<%- if @grpc_address -%>
  --grpc-address=<%=@grpc_address%> \
<%- end -%>
<%- if @grpc_advertise_address -%>
  --grpc-advertise-address=<%=@grpc_advertise_address%> \
<%- end -%>
<%- if @http_address -%>
  --http-address=<%=@http_address%> \
<%- end -%>
<%- if @http_advertise_address -%>
  --http-advertise-address=<%=@http_advertise_address%> \
<%- end -%>
<%- if @prometheus_url -%>
  --prometheus.url=<%=@prometheus_url%> \
<%- end -%>
<%- if @tsdb_path -%>
  --tsdb.path=<%=@tsdb_path%> \
<%- end -%>
<%- if @sidecar_objstore_config_file -%>
  --objstore.config-file=<%=@sidecar_objstore_config_file%> \
<%- end -%>
<%- if @reloader_config_file -%>
  --reloader.config-file=<%=@reloader_config_file%> \
<%- end -%>
<%- if @reloader_config_envsubst_file -%>
  --reloader.config-envsubst-file=<%=@reloader_config_envsubst_file%> \
<%- end -%>
<%- if @reloader_rule_dir -%>
  <% @reloader_rule_dir.each do |val| -%>
--reloader.rule-dir=<%= val %> \
  <% end -%>
<%- end -%>
<%- if @log_level -%>
  --log.level=<%=@log_level%>
<%- else -%>
  --log.level=info
<% end -%>
Restart=always

[Install]
Alias=thanos-sidecar.service
WantedBy=multi-user.target

