# templates/thanos-compact.service.erb
# This file managed by Puppet
# Template Path: <%= @module_name %>/templates/thanos-compact.service
### Managed by Puppet ###
[Unit]
Description=Prometheus Thanos Subsystem
Wants=basic.target
After=basic.target network.target

[Service]
ExecStart=/usr/bin/thanos compact \
<%- if @gcloudtrace_project -%>
  --gcloudtrace.project=<%=@gcloudtrace_project%> \
<%- end -%>
<%- if @gcloudtrace_sample_factor -%>
  --gcloudtrace.sample-factor=<%=@gcloudtrace_sample_factor%> \
<%- end -%>
<%- if @http_address -%>
  --http-address=<%=@http_address%> \
<%- end -%>
<%- if @data_dir -%>
  --data-dir=<%=@data_dir%> \
<%- end -%>
<%- if @sync_delay -%>
  --sync-delay=<%=@sync_delay%> \
<%- end -%>
<%- if @wait -%>
  --wait \
<%- end -%>
<%- if @compact_objstore_config_file -%>
  --objstore.config-file=<%=@compact_objstore_config_file%> \
<%- end -%>
<%- if @log_level -%>
  --log.level=<%=@log_level%>
<%- else -%>
  --log.level=info
<% end -%>
Restart=always

[Install]
WantedBy=multi-user.target